# ok-batch/ok-batch.cron
*  *  * * * root /var/www/ok-batch/ok-batch.sh
01 00 1 * * root /var/www/ok-batch/ok-logrotate.sh
